# salt/server-docker/etc/systemd/system/docker.service.d/overrides.conf
[Service]
EnvironmentFile=-/etc/default/docker
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// \
          $DOCKER_OPTS -H tcp://127.0.0.1:2375
